# configs/codex/systemd/codex-heartbeat.timer
[Unit]
Description=Run MyAgent Codex heartbeat every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
